# deploy/systemd/stem-frontend.service
[Unit]
Description=STEM Studio Vite Frontend
After=network.target stem-backend.service
Wants=stem-backend.service

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__PROJECT_ROOT__/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
Environment=VITE_ALLOWED_HOSTS=__PUBLIC_HOST__
ExecStart=/usr/bin/npm run preview
Restart=always
RestartSec=3
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
